# blobs/inittab
[common]
count=6

[task_1]
exec=/bin/mount
opts=-a
respawn=0

[task_2]
exec=/bin/sysctl
opts=hostname=liveiso
respawn=0

[task_3]
exec=/bin/uname
opts=-a
respawn=0

[task_4]
exec=/bin/login
opts=0
respawn=1

[task_5]
exec=/bin/login
opts=1
respawn=1

[task_6]
exec=/bin/login
opts=2
respawn=1
